# deploy/launchd/dev.sovereign-rag.frontend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dev.sovereign-rag.frontend</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/npm</string>
    <string>run</string>
    <string>dev</string>
    <string>--</string>
    <string>--host</string>
    <string>0.0.0.0</string>
    <string>--port</string>
    <string>5173</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/muhriddin/actions-work/sovereign-rag/frontend</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/bin:/bin</string>
    <key>NODE_ENV</key>
    <string>development</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/sovereign-rag-frontend.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/sovereign-rag-frontend.log</string>

  <key>ThrottleInterval</key>
  <integer>10</integer>
</dict>
</plist>
